# deploy/oracle/bankongseton-api.service
[Unit]
Description=BANKONGSETON Mobile/Web API (Flask, /api/*)
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/bankongseton/backend/api
EnvironmentFile=/opt/bankongseton/.env
ExecStart=/opt/bankongseton/venv/bin/python api_server.py
Restart=on-failure
RestartSec=3
KillMode=mixed

[Install]
WantedBy=multi-user.target
