# deploy/oracle/bankongseton-web.service
[Unit]
Description=BANKONGSETON Web Dashboard (Flask + SocketIO gunicorn)
After=network.target

[Service]
Type=notify
User=__USER__
Group=__USER__
WorkingDirectory=__APP_DIR__/backend/dashboard
EnvironmentFile=__APP_DIR__/.env
ExecStart=__APP_DIR__/venv/bin/gunicorn -c __APP_DIR__/deploy/oracle/gunicorn_web.py web_app:app
Restart=on-failure
RestartSec=3
KillMode=mixed

[Install]
WantedBy=multi-user.target
